# mufassah-tracker/mufassah-tracker.service
[Unit]
Description=Mufassah Device Tracking Service
After=network.target gpsd.service
Requires=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/mufassah-tracker
Environment="PATH=/usr/bin:/home/pi/.local/bin"
ExecStart=/usr/bin/python3 tracker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mufassah-tracker

# Resource limits
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
